# realtek/realtek/rtk-nvr/files/rtk-nvr.init
#!/bin/sh /etc/rc.common
# Copyright (C) 2017 OpenWrt.org

START=99

start() {
	export GST_OMX_CONFIG_DIR=/usr/local/etc/xdg
	[ -x /usr/bin/NVRDaemon ] && /usr/bin/NVRDaemon &
}
