# linux-4.4.18/root-rtd1295/usr/lib/opkg/info/mdadm.list
/etc/mdadm/mdadm.conf
/sbin/mdadm
/etc/init.d/mdadm-raid
