# linux-4.4.18/root-rtd1295/usr/lib/opkg/info/minidlna.list
/etc/init.d/minidlna
/etc/minidlna.conf
/usr/bin/minidlna
/etc/config/minidlna
